# systemd/sponda-revisit-reminders.timer
[Unit]
Description=Send Sponda revisit reminder emails daily at 11:00 UTC

[Timer]
OnCalendar=*-*-* 11:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
