# systemd/sponda-refresh.service
[Unit]
Description=Sponda daily data refresh (IPCA + tickers)
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/sponda/backend
Environment=DJANGO_SETTINGS_MODULE=config.settings.production
ExecStart=/opt/sponda/venv/bin/python manage.py refresh_ipca
ExecStart=/opt/sponda/venv/bin/python manage.py refresh_tickers
User=root
