# systemd/sponda-frontend.service
[Unit]
Description=Sponda Next.js Frontend
After=network.target sponda.service
Wants=sponda.service

[Service]
User=root
WorkingDirectory=/opt/sponda/frontend
EnvironmentFile=/opt/sponda/.env
Environment=DJANGO_API_URL=http://127.0.0.1:8710
Environment=PORT=3100
ExecStart=/usr/local/bin/npx next start -p 3100
Restart=always

[Install]
WantedBy=multi-user.target
